# elderberrypy/files/gentoo/init.d/ebpy-supervisord
#!/sbin/runscript
#
# Gentoo startup script for a supervisor instance
#
# description: supervisor

opts="${opts} reload"

depend(){
    need net
}

start() {
    ebegin "Starting supervisord"
    if [ ! -e ${CONF} ]; then
        eend 1 "${CONF} does not exist"
        return 1
    fi
    ${SUPERVISORD} -c $${CONFIG}
    eend $?
}

stop() {
    ebegin "Stopping supervisord"
    ${SUPERVISORCTL} -c ${CONFIG} shutdown
    eend $?
}

reload() {
    ebegin "Reloading supervisord"
    ${SUPERVISORCTL} -c ${CONFIG} reload
    eend $?
}

